# com.darren.timesheet-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label — unique identifier for the agent -->
    <key>Label</key>
    <string>com.darren.timesheet-bot</string>

    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/darren/Desktop/Projects/sharepoint_timesheet_bot/scripts/run_timesheet.sh</string>
    </array>

    <!-- Schedule: every Friday at 09:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/darren/Desktop/Projects/sharepoint_timesheet_bot</string>

    <!-- Logging (fallback — the wrapper script also logs) -->
    <key>StandardOutPath</key>
    <string>/Users/darren/Desktop/Projects/sharepoint_timesheet_bot/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/darren/Desktop/Projects/sharepoint_timesheet_bot/logs/launchd_stderr.log</string>

    <!-- Do NOT run immediately on load — wait for the schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/darren</string>
    </dict>
</dict>
</plist>
